# systemd/pyhss_gsup.service
[Unit]
Description=PyHSS GSUP Service
PartOf=pyhss.service


[Service]
Environment="PYTHONUNBUFFERED=1"
User=root
WorkingDirectory=/etc/pyhss/services/
ExecStart=python3 gsupService.py
Restart=always

[Install]
WantedBy=pyhss.service
