# systemd/pyhss_metric.service
[Unit]
Description=PyHSS Metric Service
PartOf=pyhss.service


[Service]
Environment="PYTHONUNBUFFERED=1"
User=root
WorkingDirectory=/etc/pyhss/services/
ExecStart=python3 metricService.py
Restart=always

[Install]
WantedBy=pyhss.service
