# debian/systemd/pyhss_api.service.d/distro_packaging.conf
[Service]
User=pyhss
WorkingDirectory=/var/lib/pyhss
StateDirectory=pyhss
ExecStart=
ExecStart=/opt/venvs/pyhss/bin/pyhss_api

[Install]
WantedBy=multi-user.target
